# templates/svc-template-laravel/configs/worker.conf
; Laravel workers (added to the global supervisord via /etc/supervisor/conf.d/).
; Toggle with ENABLE_QUEUE / ENABLE_SCHEDULER.
[program:queue]
command=php artisan queue:work --tries=3
directory=/var/www/html
autostart=%(ENV_ENABLE_QUEUE)s
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:scheduler]
command=sh -c "while true; do php artisan schedule:run --no-interaction; sleep 60; done"
directory=/var/www/html
autostart=%(ENV_ENABLE_SCHEDULER)s
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
