# templates/cron-template-php/crontab
# Schedule (supercronic). Times use the container TZ. fields: m h dom mon dow
*/5 * * * * php /app/job.php
