# configs/php-app/supervisord.conf
; Global default supervisord config for php-fpm app projects.
; Universal programs (php-fpm + nginx) live here; each is toggled by an env var.
; Per-project extras (queue, scheduler, vite, …) are added as their own files in
; /etc/supervisor/conf.d/ — included below.
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

[program:php-fpm]
command=php-fpm -F
autostart=%(ENV_ENABLE_PHP)s
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g 'daemon off;'
autostart=%(ENV_ENABLE_NGINX)s
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Generic cron via supercronic. Off unless ENABLE_CRON=true. A php project that
; wants it provides /etc/supervisor/crontab (COPY in its Dockerfile or mount).
; Framework templates (Laravel/Symfony) keep their OWN scheduler programs in
; /etc/supervisor/conf.d/ instead of using this one.
[program:cron]
command=supercronic -passthrough-logs /etc/supervisor/crontab
autostart=%(ENV_ENABLE_CRON)s
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[include]
files = /etc/supervisor/conf.d/*.conf
