# templates/cron-template-shell/crontab
# Schedule (supercronic). Times use the container TZ. fields: m h dom mon dow
*/5 * * * * sh /app/job.sh
